# infra/systemd/appof-web.service
[Unit]
Description=Appof Web Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/appof
Environment=HOME=/var/www
Environment=NODE_ENV=production
EnvironmentFile=/srv/appof/.env
ExecStart=/usr/bin/pnpm --filter @news/web start
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
