# infra/systemd/appof-api.service
[Unit]
Description=Appof API Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/appof/apps/api
EnvironmentFile=/srv/appof/apps/api/.env
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
